# systemd/qwen36-llama.service
[Unit]
Description=Qwen3.6-27B llama.cpp local model worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/adventor/simon
Environment=HOST=127.0.0.1
Environment=PORT=8081
Environment=CTX_SIZE=32768
Environment=BATCH_SIZE=1024
Environment=UBATCH_SIZE=256
Environment=PARALLEL=1
Environment=REASONING=off
ExecStart=/home/adventor/simon/run-qwen36-27b-server.sh
Restart=always
RestartSec=5
StandardOutput=append:/home/adventor/simon/logs/qwen36-server.log
StandardError=append:/home/adventor/simon/logs/qwen36-server.log

[Install]
WantedBy=default.target
